# gpio-shutdown.service
[Unit]
Description=GPIO Button Shutdown Listener
After=multi-user.target sysinit.target
Wants=multi-user.target

[Service]
Type=simple
EnvironmentFile=/opt/gpio-shutdown/env.conf
ExecStart=/opt/gpio-shutdown/venv/bin/python /opt/gpio-shutdown/gpio_shutdown.py

Restart=on-failure
RestartSec=5
StartLimitInterval=30s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
